# ops/com.lolwierd.rig.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lolwierd.rig</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Replace with your node/npm path: run `which npm` -->
        <string>/usr/local/bin/npm</string>
        <string>start</string>
    </array>
    <key>WorkingDirectory</key>
    <!-- Replace with the absolute path to rig/server on your machine -->
    <string>/path/to/rig/server</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/rig.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/rig.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <!-- Ensure the directory containing `node` and `pi` is included -->
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
