# ops/rig.service
[Unit]
Description=Rig Coding Agent Platform
After=network.target

[Service]
Type=simple
# Replace with your username
User=YOUR_USER
# Replace with the absolute path to rig/server on your machine
WorkingDirectory=/path/to/rig/server
# Replace with your npm path: run `which npm`
ExecStart=/usr/local/bin/npm start
Restart=always
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
